# webalizer.crontab
#
# cron job for the webalizer
#
05 * * * *	stats	/usr/bin/webalizer.cron
