# webalizer.cron
#!/bin/sh

[ -f /etc/sysconfig/webalizer ] && . /etc/sysconfig/webalizer

if [ -n "$ENABLE" ] && [ "$ENABLE" = "YES" ]; then
	LC_ALL=$LANG
	export LANG
	export LC_ALL
	
	if [ -f /etc/webalizer.conf ]; then
		webalizer -c /etc/webalizer.conf
	fi
	
	for i in /etc/webalizer/*.conf; do webalizer -c $i; done
fi
