# webalizer.crontab
#
# cron job for the webalizer
#
* * * * *	stats	/usr/bin/webalizer.cron
